# rootfs_master/usr/service/kube-controller-manager/run
#!/bin/bash

exec 2>&1
exec \
  envdir ./env \
    sh -c 'exec /usr/local/bin/controller-manager   \
 	--master=127.0.0.1:8080 \
	--node-monitor-grace-period=10s \
	--pod-eviction-timeout=20s \
	--service-account-private-key-file=/srv/kubernetes/server.key \
	--root-ca-file=/srv/kubernetes/ca.crt \
	--service-cluster-ip-range="10.253.0.0/16" \
 	--logtostderr=true \
	--v 1 \
	'

