# rootfs_node/usr/service/kube-proxy/run
#!/bin/bash

exec 2>&1
exec \
  envdir ./env \
    sh -c 'exec /usr/local/bin/proxy  \
        --master=http://${MASTER_IP}:8080 \
	--masquerade-all=true \
        --v=3 \
        '

